# bus/dbus.socket.in
[Socket]
ListenStream=@DBUS_SYSTEM_SOCKET@

[Install]
WantedBy=sockets.target
